# com.ora.weekly-summary.plist

<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ora.weekly-summary</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/air/ea_assistant/generate_weekly_summary.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
        <key>Weekday</key>
        <integer>1</integer> <!-- Monday -->
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/air/Library/Logs/ora/weekly_summary.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/air/Library/Logs/ora/weekly_summary.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
